# deploy/lucid-bridge.service
[Unit]
Description=Lucid bridge (relay + ingest + OCR)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Replace __USER__ + __DIR__ at install time (see README-aws.md).
User=__USER__
WorkingDirectory=__DIR__
ExecStart=/usr/bin/node index.js
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
